# service/uwsgi/run
#!/bin/sh

exec 2>&1
realpath "$0"
set -ex
exec uwsgi --ini web2py.ini
